# config/supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:runserver]
directory=/code
command=python manage.py runserver 0.0.0.0:8000
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/django.log
stderr_logfile=/var/log/supervisor/django_error.log
environment=DJANGO_SETTINGS_MODULE="core.settings",PYTHONUNBUFFERED="1"

[program:daphne]
directory=/code
command=daphne -b 0.0.0.0 -p 8001 core.asgi:application
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/daphne.log
stderr_logfile=/var/log/supervisor/daphne_error.log
environment=DJANGO_SETTINGS_MODULE="core.settings",PYTHONUNBUFFERED="1"
